# jetson_scripts/rugbycodex-cleanup.service
[Unit]
Description=RugbyCodex Media Cleanup Worker
After=network.target

[Service]
Type=simple
User=scribe
WorkingDirectory=/home/scribe/rugbycodex/jetson_scripts/cleanup
EnvironmentFile=/home/scribe/rugbycodex/jetson_scripts/.env
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 cleanup_worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
